# scripts/akmon.service
[Unit]
Description=Akmon — eBPF-based host monitoring for AI agent processes
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root

ExecStart=/usr/local/bin/akmon \
    --ui :9090 \
    --config /etc/akmon/config.yaml \
    --bpf-obj   /usr/lib/akmon/monitor.bpf.o \
    --behavioral-templates /etc/akmon/behavioral-templates \
    --nuclei-templates /etc/akmon/nuclei-templates \
    --log-level info \
    --output    /var/log/akmon/events.json \
    --compact \
    --compact-log /var/log/akmon/chains.json


# View live: journalctl -u akmon -f
# View file: tail -f /var/log/akmon/monitor.log
StandardOutput=journal
StandardError=append:/var/log/akmon/monitor.log

Restart=on-failure
RestartSec=3

# Graceful stop: systemd (PID 1) sends SIGTERM — our LSM hook allows PID 1.
KillSignal=SIGTERM
KillMode=control-group
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
